# dot_config/systemd/user/trash-empty@.service
[Unit]
Description=Delete trashed files older than %i days

[Service]
Type=exec
ExecStart=trash-empty %i
