# xdg/dot-config/systemd/user/trash-empty.timer
[Unit]
Description=Empty the trash

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
