# xdg/dot-config/systemd/user/trash-empty.service
[Unit]
Description=Empty the trash

[Service]
ExecStart=trash-empty 7
