# helix/dot-config/profile.d/helix.sh
# Use Helix as the default editor
export EDITOR='hx'
export VISUAL="$EDITOR"
