# sh/dot-config/profile.d/02-man.bash
# Configure man
export MANLESS=$'manual $MAN_PN lines %lt-%lb/%L%t'
export MANROFFOPT='-c'
