# rclone/dot-config/systemd/user/rclone.target
[Unit]
Description=Rclone jobs

[Install]
WantedBy=timers.target
