# ov/dot-config/profile.d/ov.sh
# Use ov as the default pager
export PAGER='ov'

# Use ov as the pager for man
export MANPAGER="ov --section-delimiter '^[^\s]' --section-header"
